# examples/integrations/systemd/artifact-watchdog.timer
[Unit]
Description=Run artifact-watchdog every morning

[Timer]
OnCalendar=*-*-* 06:30:00
Persistent=true
Unit=artifact-watchdog.service

[Install]
WantedBy=timers.target
